# LaunchAgents/com.yoshihide.run_sync.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>EnvironmentVariables</key>
	<dict>
		<key>PATH</key>
		<string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
	</dict>
	<key>Label</key>
	<string>com.yoshihide.run_sync</string>
	<key>ProgramArguments</key>
	<array>
		<string>/Users/yoshihide/my-projects/spotify-playlist-tools/sync.sh</string>
	</array>
	<key>RunAtLoad</key>
	<false/>
	<key>StandardErrorPath</key>
	<string>/Users/yoshihide/Library/Logs/com.yoshihide.run_sync.err</string>
	<key>StandardOutPath</key>
	<string>/Users/yoshihide/Library/Logs/com.yoshihide.run_sync.out</string>
	<key>StartCalendarInterval</key>
	<dict>
		<key>Hour</key>
		<integer>12</integer>
		<key>Minute</key>
		<integer>0</integer>
	</dict>
</dict>
</plist>
